# squashfs-root/etc/uci-defaults/luci-theme-lantiq-default
#!/bin/sh
uci batch <<-EOF
	set luci.main.mediaurlbase=/luci-static/lantiq
	commit luci
EOF
