# squashfs-root/etc/uci-defaults/luci-theme-lantiq
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Lantiq=/luci-static/lantiq
	commit luci
EOF
	
