# squashfs-root/usr/lib/opkg/info/gpon-omci-api.list
/opt/lantiq/bin/omci_usock_server
/etc/init.d/usock.sh
