# squashfs-root/usr/lib/opkg/info/gpon-base-files.list
/etc/uci-defaults/timeserver
/etc/init.d/show_version.sh
/opt/lantiq/image_version
/etc/uci-defaults/firewall
/opt/lantiq/bin/goi_store_uboot.sh
/opt/lantiq/bin/gpio_setup.sh
/opt/lantiq/bin/falcon_gpio_config.sh
/etc/init.d/fw_env.sh
/lib/falcon.sh
/etc/goiupgrade.conf
/etc/config/sfp_pins
/etc/init.d/pin_cfg.sh
/etc/config/gpon
/opt/lantiq/bin/whatz.sh
/etc/diag.arch.sh
/etc/config/omci
/etc/init.d/clkconfig.sh
/opt/lantiq/bin/what.sh
/opt/lantiq/bin/tohex
/etc/profile.arch
/opt/lantiq/bin/show_version.sh
/etc/init.d/recover_goi.sh
/etc/uci-defaults/hostname
/etc/uci-defaults/netdev-cfg
/etc/init.d/unaligned_action.sh
