# squashfs-root/usr/lib/opkg/info/gpon-omci-onu.list
/opt/lantiq/bin/omci_simulate.sh
/opt/lantiq/bin/omci_hook.sh
/opt/lantiq/bin/config/onu_setup_no_vlan.txt
/opt/lantiq/bin/omci_pipe.sh
/opt/lantiq/bin/omci_simulate
/etc/mibs/sfp_8q_mib.ini
/opt/lantiq/bin/config/local_bridge.txt
/etc/init.d/omcid.sh
/etc/mibs/sfu_voice.ini
/etc/mibs/1ge_sfu_voice.ini
/etc/mibs/occam.ini
/etc/mibs/sfp.ini
/etc/mibs/sfp_alu.ini
/opt/lantiq/bin/config/onu_setup_vlan.txt
/etc/mibs/alu.ini
/etc/mibs/sfu_voice_8q_mib.ini
/etc/mibs/mdu.ini
/opt/lantiq/bin/omcid
/opt/lantiq/bin/host2lan.sh
/opt/lantiq/bin/local_bridge.sh
