# squashfs-root/usr/lib/opkg/info/gpon-sfp-eeprom-linux.list
/opt/lantiq/bin/sfp_i2c
/etc/init.d/sfp_eeprom.sh
/etc/config/sfp_eeprom
/etc/uci-defaults/uci_disable_dying_gasp.sh
