# squashfs-root/usr/lib/opkg/info/gpon-ocal.list
/opt/lantiq/bin/ocal
/opt/lantiq/bin/ocal_pipe.sh
/etc/init.d/load_ocal.sh
