# squashfs-root/usr/lib/opkg/info/gpon-dti-agent.list
/etc/init.d/load_dti_agent.sh
/opt/lantiq/bin/gpon_dti_agent
