# squashfs-root/usr/lib/opkg/info/gpon-optic-drv.list
/opt/lantiq/bin/optic
/etc/init.d/optic.sh
/opt/lantiq/bin/otop
/etc/hotplug.d/gpon/10-optic
