# squashfs-root/usr/lib/opkg/info/ipwatchd.list
/usr/share/ipwatchd/ipwatchd.script
/etc/init.d/ipwatchd
/usr/sbin/ipwatchd
/etc/ipwatchd.conf
